# docker/supervisord.conf
[unix_http_server]
file=/run/supervisord.sock

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/supervisord.sock

[program:web-ui]
command=node /app/web/build/index.js
directory=/app/web
autostart=true
autorestart=true
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME="/root",PORT="8080",ORIGIN="https://fabric.aws.xin",BODY_SIZE_LIMIT="100M"

[program:fabric-api]
command=/bin/sh -c "exec /usr/local/bin/fabric --serve --address :8180 < /dev/null"
autostart=true
autorestart=true
startsecs=5
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME="/root",FABRIC_CONFIG_DIR="/root/.config/fabric"

[program:mcp-server]
command=python3 /app/mcp_server.py
autostart=true
autorestart=true
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME="/root",MCP_PORT="8181"
